# Pi-s/RaspberryPi监控/RaspberryPi搭内网摄像头/camera-stream.service
[Unit]
Description=Pi Camera on-demand MJPEG stream
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/bin/python3 /home/pi/camera_stream_on_demand.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
